# wifi-heartbeat.timer
[Unit]
Description=Run Wi-Fi Heartbeat every 30s

[Timer]
OnBootSec=15s
OnUnitActiveSec=30s
AccuracySec=1s
Unit=wifi-heartbeat.service

[Install]
WantedBy=timers.target